# templates/update-motd.d/99-server-tools
#!/bin/sh

echo
server-tools version:check nodejs
echo