# templates/update-motd.d/99-stool
#!/bin/sh

stool version:check nodejs
stool version:check composer
stool version:check nvm