# templates/update-motd.d/99-server-tools
#!/bin/sh

server-tools version:check nodejs
server-tools version:check composer
server-tools version:check nvm