# templates/update-motd.d/99-server-tools
#!/bin/sh

server-tools version:check nodejs